# systemd/assignment-sensor.service
[Unit]
Description=Mock Sensor Logger Service
After=network.target

[Service]
ExecStart=/usr/local/bin/sensor_logger --interval 5
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
